# deploy/aster-support-bot.service
# ============================================================================
# Aster Support Bot - systemd Service Unit
# 阿里云 Telegram 支持机器人 - systemd 服务单元
# 
# Security Features | 安全特性:
# - Non-root execution | 非root执行
# - KMS secret integration | KMS密钥集成
# - Minimal privileges | 最小权限
# - Sandboxing enabled | 沙箱隔离启用
# ============================================================================

[Unit]
Description=Aster Telegram Support Bot
Description[zh_CN]=阿里云 Telegram 客户支持机器人
Documentation=https://github.com/ross-32/Telegram-Support-Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# User and Group | 用户和组
# IMPORTANT: Create user before starting service | 重要：启动服务前创建用户
User=asterbot
Group=asterbot

# Working Directory | 工作目录
WorkingDirectory=/opt/aster-support-bot

# KMS Configuration | KMS 配置
# These are NOT secrets, just metadata to locate the secret | 这些不是密钥，仅是定位密钥的元数据
Environment="KMS_SECRET_NAME=aster/telegram_support_bot/config"
Environment="KMS_REGION_ID=cn-hangzhou"
Environment="KMS_VERSION_ID="

# Environment File | 环境变量文件
# This file will be populated by ExecStartPre script | 此文件将由 ExecStartPre 脚本填充
EnvironmentFile=/etc/aster-support-bot/aster-support-bot.env

# Pre-start: Fetch secrets from KMS | 启动前：从 KMS 获取密钥
# This script runs as root to write env file, then service drops to asterbot user
# 此脚本以 root 运行以写入 env 文件，然后服务降权到 asterbot 用户
ExecStartPre=/usr/local/bin/aster-support-bot-fetch-secrets.sh

# Main Process | 主进程
ExecStart=/opt/aster-support-bot/venv/bin/python /opt/aster-support-bot/bot.py

# Restart Policy | 重启策略
Restart=always
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=60

# Logging (journald) | 日志（journald）
# Logs can be viewed with: journalctl -u aster-support-bot -f
# 日志可通过以下命令查看: journalctl -u aster-support-bot -f

# Optional: Redirect to files | 可选：重定向到文件
# Uncomment to enable file logging | 取消注释以启用文件日志
# StandardOutput=append:/var/log/aster-support-bot/bot.log
# StandardError=append:/var/log/aster-support-bot/bot.error.log

# Security Hardening | 安全加固
# Prevents privilege escalation | 防止权限提升
NoNewPrivileges=true

# Private /tmp | 私有 /tmp
PrivateTmp=true

# Read-only root filesystem, except specific paths | 只读根文件系统，除特定路径外
ProtectSystem=strict
ReadWritePaths=/opt/aster-support-bot
ReadWritePaths=/var/log/aster-support-bot

# Protect /home | 保护 /home
ProtectHome=true

# Restrict file creation permissions | 限制文件创建权限
UMask=0077

# Additional sandboxing (optional, uncomment if compatible) | 额外沙箱（可选，如兼容请取消注释）
# ProtectKernelTunables=true
# ProtectKernelModules=true
# ProtectControlGroups=true
# RestrictRealtime=true
# RestrictNamespaces=true
# LockPersonality=true
# MemoryDenyWriteExecute=true
# RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
